# examples/service/amex-ynab-import.timer
[Timer]
# Schedule the service to run at noon and midnight every day
OnCalendar=*-*-* 12,0:00:00